# ddns.timer
[Unit]
Description=ddns timer

[Timer]
OnUnitActiveSec=200s
OnBootSec=200s

[Install]
WantedBy=timers.target
